# kasse.service
[Unit]
Description=Kasse – Ausgaben Tracker API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/kasse/backend
ExecStart=/opt/kasse/venv/bin/gunicorn app:app \
    --bind 127.0.0.1:5000 \
    --workers 2 \
    --timeout 30 \
    --access-logfile /var/log/kasse/access.log \
    --error-logfile /var/log/kasse/error.log
Environment="DB_PATH=/data/kasse/kasse.db"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
